# docs/deploy/gunicorn/service_bot.socket
[Unit]
Description=gunicorn socket

[Socket]
ListenStream=/run/service_bot.sock

[Install]
WantedBy=sockets.target
